# systemd/casa-mounts.service
[Unit]
Description=Casa NAS mount readiness — ping Unraid, mount SMB shares, verify
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=casaroot
ExecStart=/home/casaroot/apps/mounts_ready.sh
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
