# systemd/casa-planetexpress.service.template
[Unit]
Description=Planet Express — Docker Compose Homelab Sysadmin Agent (Farnsworth)
Documentation=https://github.com/sovereignalmida/planet-express
After=network-online.target docker.service
Wants=network-online.target
# Don't start until Docker is up
Requires=docker.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=$RUN_USER
Group=$RUN_GROUP
WorkingDirectory=$INSTALL_DIR

# Load secrets (Telegram bot token, LLM API key) from a separate file
EnvironmentFile=/etc/planetexpress.env

Environment=CASA_CONFIG=$CONFIG_FILE
Environment=CASA_STATE_DIR=$INSTALL_DIR/state
Environment=CASA_LOG_DIR=$INSTALL_DIR/logs
Environment=PYTHONUNBUFFERED=1

ExecStart=$INSTALL_DIR/venv/bin/python casa_farnsworth.py

# Restart on failure with backoff
Restart=on-failure
RestartSec=15s

# Logging — journald picks this up via stdout
StandardOutput=journal
StandardError=journal
SyslogIdentifier=casa-planetexpress

[Install]
WantedBy=multi-user.target
